# nfc-music-player/nfc-music.service
[Unit]
Description=NFC Music Player
After=multi-user.target
Wants=sound.target

[Service]
Type=simple
User=cody
Group=cody
WorkingDirectory=/home/cody/nfc-music-player
ExecStart=/usr/bin/python3 /home/cody/nfc-music-player/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target